# rc.local

# Paste this line in /etc/rc.local file and create needed directories in /home/pi

sudo python /home/pi/projects/OLEDSystemStats/stats.py &